# init/runit/log/run
#!/bin/sh
exec chpst -u omnisearch:omnisearch svlogd -tt /var/log/omnisearch