# init/s6/log/run
#!/bin/execlineb -P
exec -c chpst -u omnisearch:omnisearch s6-log -d 3 /var/log/omnisearch